# files/etc/uci-defaults/99-apply-confirm
#!/bin/sh
mkdir -p /etc/apply-confirm/pending
exit 0
